# nightupdate/etc/systemd/system/arch-nightly-update.timer
[Unit]
Description=Nightly Arch Linux Update Timer

[Timer]
OnCalendar=*-*-* 03:30
Persistent=true

[Install]
WantedBy=timers.target

